# extra/upstart/rlog.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

script
	cd /srv/rlog
	exec /srv/rlog/rlog
end script
